# roles/acme_client/files/dehydrated.service
[Unit]
Description=Renew ACME-derived certificates with dehydrated

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dehydrated --cron
